# sites/dallas-site/supervisor/filling-line.conf
; ============================================================================
; dallas-filling-line (gw-fill-01) — FULLY SELF-CONTAINED edge device.
;
; One container = one edge gateway; every edgecommons component runs as a supervised
; process, "bare-metal style", against the container's OWN in-process EMQX broker on
; localhost:1883. The field sims live IN the container too, so the adapters read them
; over loopback (no cross-container OPC UA endpoint-rewrite needed anymore).
;
; START ORDER (priority, lower first) + READINESS GATING (wait-for-tcp):
;   10  emqx            local device bus (localhost:1883)
;   20  opcua-sim       asyncua  -> binds opc.tcp://localhost:4840/
;   20  modbus-sim      pymodbus -> binds 0.0.0.0:5020
;   25  config-component serves the line's hierarchical catalog
;   28  uns-bridge      waits 1883; relays localhost:1883 -> dallas-site:1883 (retries site)
;   30  opcua-adapter   waits 1883 + 4840, then connects the sim
;   30  modbus-adapter  waits 1883 + 5020, then connects the sim
;   40  telemetry       waits 1883; subscribes data -> Lua -> Parquet /out/archive
;   50  file-replicator waits 1883; watches /out/archive -> /mnt/replicated (host bind)
;
; This same image runs the packaging line too (packaging-line.conf) — it simply does not
; start the sims / telemetry / file-replicator.
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; Local control socket so `supervisorctl status` works inside the container
; (introspect which of the many co-located processes is up/down/FATAL).
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:emqx]
command=/usr/bin/emqx foreground
user=emqx
environment=HOME="/var/lib/emqx",EMQX_ALLOW_ANONYMOUS="true"
priority=10
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Dallas filling-line field sim: ONE process serves OPC UA :4840 + Modbus :5020 from a shared
# scenario (SIMULATOR_CONFIGURATION.md "Line 1"), replacing the two generic adapter smoke sims.
[program:field-sim]
command=/opt/simenv/bin/python /opt/sims/dallas_filling_sim.py
priority=20
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:config-component]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- config-component --platform HOST --transport MQTT /config/config-component-messaging.json -c FILE /config/config-component-config.json -t gw-fill-01
priority=25
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:opcua-adapter]
command=/usr/local/bin/wait-for-tcp localhost:1883 localhost:4840 -- /bin/bash -lc 'sleep 6; exec java --enable-native-access=ALL-UNNAMED -jar /app/opcua/app.jar --platform HOST --transport MQTT /config/opcua-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=30
autorestart=true
startsecs=8
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:modbus-adapter]
command=/usr/local/bin/wait-for-tcp localhost:1883 localhost:5020 -- /bin/bash -lc 'sleep 6; exec /opt/pyenv/bin/python /app/modbus/main.py --platform HOST --transport MQTT /config/modbus-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
directory=/app/modbus
priority=30
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:telemetry-processor]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 6; exec telemetry-processor --platform HOST --transport MQTT /config/telemetry-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=40
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:file-replicator]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 6; exec file-replicator --platform HOST --transport MQTT /config/filerep-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=50
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uns-bridge]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 3; exec uns-bridge --platform HOST --transport MQTT /config/bridge-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=28
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
