# sites/dallas-site/supervisor/packaging-line.conf
; ============================================================================
; dallas-packaging-line (gw-pack-01) — edge device with EXTERNAL field sources.
;
; Same edge-node image as the filling line, but this supervisord conf starts a REDUCED
; process set: EMQX + the two adapters + the bridge. There are NO field sims, NO telemetry
; processor and NO file-replicator here — the packaging adapters read REAL plant PLCs on
; the LAN (KepWare OPC UA + a host Modbus/TCP sim), so the container is not self-contained.
;
; The ConfigComponent renders the hierarchical catalog template (host:port / endpoint injected
; from the container environment — set by the compose service from .env) into /run/config
; before it serves CONFIG_COMPONENT requests. The adapters then use the standard provider path.
;
; START ORDER (priority) + READINESS GATING:
;   10  emqx            local device bus (localhost:1883)
;   25  config-component render/load the line's hierarchical catalog
;   28  uns-bridge      wait 1883; relay localhost:1883 -> dallas-site:1883 (retries site)
;   30  opcua-adapter   wait 1883; load CONFIG_COMPONENT; connect ${KEPWARE_ENDPOINT}
;   30  modbus-adapter  wait 1883; load CONFIG_COMPONENT; connect ${HOST_MODBUS}
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; Local control socket so `supervisorctl status` works inside the container.
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:emqx]
command=/usr/bin/emqx foreground
user=emqx
environment=HOME="/var/lib/emqx",EMQX_ALLOW_ANONYMOUS="true"
priority=10
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:config-component]
command=/bin/bash -c 'render-packaging-catalog /config/config-catalog.tmpl.json /run/config/config-catalog.json "${KEPWARE_ENDPOINT:-opc.tcp://192.168.1.180:49320}" "${HOST_MODBUS:-192.168.1.224:5020}" && exec /usr/local/bin/wait-for-tcp localhost:1883 -- config-component --platform HOST --transport MQTT /config/config-component-messaging.json -c FILE /config/config-component-config.json -t gw-pack-01'
priority=25
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:opcua-adapter]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 6; exec java --enable-native-access=ALL-UNNAMED -jar /app/opcua/app.jar --platform HOST --transport MQTT /config/opcua-messaging.json -c CONFIG_COMPONENT -t gw-pack-01'
priority=30
autorestart=true
startsecs=8
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:modbus-adapter]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 6; exec /opt/pyenv/bin/python /app/modbus/main.py --platform HOST --transport MQTT /config/modbus-messaging.json -c CONFIG_COMPONENT -t gw-pack-01'
directory=/app/modbus
priority=30
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uns-bridge]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 3; exec uns-bridge --platform HOST --transport MQTT /config/bridge-messaging.json -c CONFIG_COMPONENT -t gw-pack-01'
priority=28
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
