# sites/dallas-site/supervisor/packaging-line.conf
; ============================================================================
; dallas-packaging-line (gw-pack-01) — edge device with EXTERNAL field sources.
;
; Same edge-node image as the filling line, but this supervisord conf starts a REDUCED
; process set: EMQX + the two adapters + the bridge. There are NO field sims, NO telemetry
; processor and NO file-replicator here — the packaging adapters read REAL plant PLCs on
; the LAN (KepWare OPC UA + a host Modbus/TCP sim), so the container is not self-contained.
;
; Each adapter renders its config template (host:port / endpoint / creds injected from the
; container environment — set by the compose service from .env) into /run/config BEFORE it
; execs, then gates on the local EMQX with wait-for-tcp. The render is a baked script so the
; supervisord command line stays free of '%' (a supervisord expansion sigil).
;
; START ORDER (priority) + READINESS GATING:
;   10  emqx            local device bus (localhost:1883)
;   30  opcua-adapter   render KepWare cfg; wait 1883; connect ${KEPWARE_ENDPOINT}
;   30  modbus-adapter  render host cfg;    wait 1883; connect ${HOST_MODBUS}
;   60  uns-bridge      wait 1883; relay localhost:1883 -> dallas-site:1883 (retries site)
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; Local control socket so `supervisorctl status` works inside the container.
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:emqx]
command=/usr/bin/emqx foreground
user=emqx
environment=HOME="/var/lib/emqx",EMQX_ALLOW_ANONYMOUS="true"
priority=10
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:opcua-adapter]
command=/bin/bash -c 'render-opcua-config /config/opcua-config.tmpl.json /run/config/opcua-config.json "${KEPWARE_ENDPOINT:-opc.tcp://192.168.1.180:49320}" "${KEPWARE_USER:-testuser}" "${KEPWARE_PASS:-CHANGE_ME}" && exec /usr/local/bin/wait-for-tcp localhost:1883 -- java --enable-native-access=ALL-UNNAMED -jar /app/opcua/app.jar --platform HOST --transport MQTT /config/opcua-messaging.json -c FILE /run/config/opcua-config.json -t gw-pack-01'
priority=30
autorestart=true
startsecs=8
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:modbus-adapter]
command=/bin/bash -c 'render-modbus-config /config/modbus-config.tmpl.json /run/config/modbus-config.json "${HOST_MODBUS:-192.168.1.224:5020}" && cd /app/modbus && exec /usr/local/bin/wait-for-tcp localhost:1883 -- /opt/pyenv/bin/python /app/modbus/main.py --platform HOST --transport MQTT /config/modbus-messaging.json -c FILE /run/config/modbus-config.json -t gw-pack-01'
directory=/app/modbus
priority=30
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uns-bridge]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- uns-bridge --config /config/bridge-config.json --thing gw-pack-01
priority=60
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
